# armbian/userpatches/extensions/airwaves-os/config/templates/systemd-airwaves-growfs.service
[Unit]
Description=Airwaves OS - grow root filesystem to fill disk
Documentation=https://airwavesos.com
After=local-fs.target
Before=docker.service airwaves-containers.service airwaves-init.service
ConditionPathExists=/opt/airwaves/scripts/airwaves-growfs

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/airwaves/scripts/airwaves-growfs
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
