# armbian/userpatches/extensions/airwaves-os/config/templates/systemd-airwaves-init.service
[Unit]
Description=Airwaves OS Initialization
Documentation=https://airwavesos.com
# airwaves-init runs docker load / network create / compose up, so the daemon
# must already be up. Order AFTER docker and BEFORE the steady-state container
# service (which does the same compose up on later boots). No network-online dep:
# first boot must work offline by loading the pre-baked image tarballs.
After=local-fs.target docker.service
Requires=docker.service
Before=airwaves-containers.service
ConditionPathExists=/opt/airwaves/.needs-first-run

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/airwaves/scripts/airwaves-init
TimeoutStartSec=5min
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
