# armbian/armbian-userpatches/extensions/airwaves-os/config/templates/systemd-airwaves-manager.service.template
[Unit]
Description=Airwaves Manager
Documentation=https://airframes.io
After=network.target local-fs.target

[Service]
ExecStart=/opt/airwaves/bin/manager

Type=notify
Delegate=yes
KillMode=process
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
